# monitoring/schedules/auto-check.cron
# 自動チェックスケジュール
# AI開発者の虚偽報告を防ぐための多層防御

# 10分ごと - 簡易チェック（ファイル数、基本構造）
*/10 * * * * node /mnt/c/Users/user/ai-template-lab-restructured/monitoring/scripts/quick-check.js

# 30分ごと - 整合性チェック（チェックサム、権限）
*/30 * * * * node /mnt/c/Users/user/ai-template-lab-restructured/monitoring/scripts/integrity-check.js

# 1時間ごと - 深層スキャン（コード品質、ルール遵守）
0 * * * * node /mnt/c/Users/user/ai-template-lab-restructured/monitoring/scripts/deep-scan.js

# 日次 - 完全監査（全項目の総合チェック）
0 0 * * * node /mnt/c/Users/user/ai-template-lab-restructured/monitoring/scripts/full-audit.js

# 週次 - テンプレート改善レビュー
0 0 * * 0 node /mnt/c/Users/user/ai-template-lab-restructured/monitoring/scripts/weekly-review.js